# templates/vault.service.erb
[Unit]
Description=Vault secret management tool
Requires=network-online.target
After=network-online.target
[Service]
User=<%= scope['vault::user'] %>
Group=<%= scope['vault::group'] %>
PIDFile=/var/run/vault/vault.pid
ExecStart=<%= scope['vault::bin_dir'] %>/vault server -config=<%= scope['vault::config_dir'] %>/config.json -log-level=debug
ExecReload=/bin/kill -HUP
KillMode=process
KillSignal=SIGTERM
Restart=on-failure
RestartSec=42s
LimitMEMLOCK=infinity
[Install]
WantedBy=multi-user.target
